# resources/_indexer.conf
start on runlevel [2345]
stop on runlevel [!2345]

exec bash /usr/local/packages/atrax_indexer/run.sh {{ crawl_job }}
